# site-cookbooks/redmine/templates/default/upstart.conf.erb
start on filesystem
respawn
<% if @kill_signal %>
kill signal <%= @kill_signal %>
<% end %>
exec <%= @command %>